# assets/systemd/pstate-frequency.service
[Unit]
Description=Program to manage CPU frequencies

[Service]
Type=simple
ExecStart=pstate-frequency -S -p 0

[Install]
WantedBy=multi-user.target
